# master-crontab.cron
#!/bin/sh

# Copyright (C) 2013 Jonathan Geisler
# Replace touche with the correct user folder


if [ `ls $HOME/active-contests | wc -l` -gt 0 ] ; then
	for file in $HOME/active-contests/*
	do
		dir=`basename $file`
		cd "$HOME/public_html/$dir/judge"
		./cronScript.php >> ./errorLog.txt &

	done
fi
