# master-crontab.cron
#!/bin/sh

# Copyright (C) 2013 Jonathan Geisler
# Replace touche with the correct user folder

if [ `ls /home/mschmock/active-contests | wc -l` -gt 0 ] ; then
	for file in /home/mschmock/active-contests/*
	do
		dir=`basename $file`
		cd "/home/mschmock/public_html/$dir/judge"
		./cronScript.php > ./errorLog.txt &
	done
fi
